# deploy/systemd/knowledge-api.service
[Unit]
Description=MinerU Knowledge API (NestJS)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__PLATFORM_DIR__/api
EnvironmentFile=__PLATFORM_DIR__/api/.env
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledge-api

# Upload / parse / index có thể lâu
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
