# deploy/systemd/mineru-api.service
[Unit]
Description=MinerU Parse API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__MINERU_DIR__
Environment=HF_HUB_DISABLE_SYMLINKS=1
Environment=MINERU_DEVICE_MODE=cpu
Environment=CUDA_VISIBLE_DEVICES=
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=ATEN_CPU_CAPABILITY=DEFAULT
# /bin/bash — tránh systemd 203/EXEC (thiếu +x hoặc CRLF trên run-mineru.sh)
ExecStart=/bin/bash __MINERU_DIR__/run-mineru.sh api
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mineru-api

# Parse PDF có thể chạy lâu, tránh kill giữa chừng
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
